# systemd/claude-memory-dreaming.service
[Unit]
Description=Claude Command Center — Memory dreaming sweep

[Service]
Type=oneshot
WorkingDirectory=%h/JjayFiles/claude-command-center
ExecStart=%h/JjayFiles/claude-command-center/.venv/bin/python memory/dreaming.py
Environment=ANTHROPIC_API_KEY=

StandardOutput=append:%h/JjayFiles/claude-command-center/logs/dreaming.log
StandardError=append:%h/JjayFiles/claude-command-center/logs/dreaming.log
