# systemd/claude-dashboard.service
[Unit]
Description=Claude Command Center — Dashboard
After=network-online.target claude-main.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/JjayFiles/claude-command-center
ExecStart=%h/JjayFiles/claude-command-center/.venv/bin/python -m uvicorn dashboard.app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

Environment=ANTHROPIC_API_KEY=

StandardOutput=append:%h/JjayFiles/claude-command-center/logs/dashboard.log
StandardError=append:%h/JjayFiles/claude-command-center/logs/dashboard.log

[Install]
WantedBy=default.target
